# roots/local/files/etc/update-motd.d/20-citraweb-local
#!/bin/sh
echo "###############################################################################"
echo "Hello! Welcome to the Citra-web local setup."
echo "For your debugging convenience, a server is not configured by default."
echo "The following commands are installed:"
echo "  - devserver: Runs a test server. Visit http://localhost:8800 to see it."
echo "  - migrate: Runs migrations on the Django application."
echo "  - makemigrations: Creates automatic migrations on Django models."
echo "  - collectstatic: Collects static files."
echo "###############################################################################"
